# ui-dashboard/systemd/leaf-can-dashboard.service
[Unit]
Description=Nissan Leaf CAN Dashboard (LVGL)
Documentation=https://github.com/Hercules-Tele-com/test-lvgl-cross-compile
After=network.target telemetry-logger-unified.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/emboo/Projects/test-lvgl-cross-compile/ui-dashboard/build
ExecStart=/home/emboo/Projects/test-lvgl-cross-compile/ui-dashboard/build/leaf-can-dashboard
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=256M
CPUQuota=50%

# Environment - Remove DISPLAY for framebuffer mode
# Environment="DISPLAY=:0"

[Install]
WantedBy=multi-user.target
