# kiosk/dashboard-kiosk.service
[Unit]
Description=Leaf Dashboard Kiosk Mode
Documentation=https://github.com/Hercules-Tele-com/test-lvgl-cross-compile
After=web-dashboard.service graphical.target network.target
Wants=web-dashboard.service graphical.target

[Service]
Type=simple
User=emboo
Group=emboo
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/emboo/.Xauthority
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/chromium \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-features=TranslateUI \
    --disable-component-update \
    --check-for-update-interval=31536000 \
    --disable-backgrounding-occluded-windows \
    --disable-breakpad \
    --disable-component-extensions-with-background-pages \
    --disable-extensions \
    --disable-features=Translate \
    --disable-ipc-flooding-protection \
    --disable-renderer-backgrounding \
    --force-device-scale-factor=1 \
    --app=http://localhost:8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical.target
