# telemetry-system/systemd/cloud-sync.service
[Unit]
Description=InfluxDB Cloud Sync Service
Documentation=https://github.com/your-repo/leaf-can-network
After=network-online.target telemetry-logger.service
Wants=network-online.target

[Service]
Type=simple
User=emboo
Group=emboo
WorkingDirectory=/home/emboo/Projects/test-lvgl-cross-compile/telemetry-system/services/cloud-sync
EnvironmentFile=/home/emboo/Projects/test-lvgl-cross-compile/telemetry-system/config/influxdb-local.env
EnvironmentFile=-/home/emboo/Projects/test-lvgl-cross-compile/telemetry-system/config/influxdb-cloud.env
Environment=PYTHONUNBUFFERED=1
Environment=SYNC_INTERVAL_SECONDS=300
Environment=NETWORK_CHECK_URL=https://1.1.1.1
ExecStart=/usr/bin/python3 /home/emboo/Projects/test-lvgl-cross-compile/telemetry-system/services/cloud-sync/cloud_sync.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
